# common/gratia-security/scripts/voms-server.cron
42 1 * * * /usr/share/gratia/voms-server.sh  >/dev/null 2>&1
